# scripts/poolassist.service
[Unit]
Description=PoolCue Vision Assist
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/PoolCue-Assist
ExecStart=/usr/bin/python3 /home/pi/PoolCue-Assist/src/main.py
Restart=always
RestartSec=5
# Pass Duke LiteLLM API key — replace with your actual key
Environment=DUKE_API_KEY=your_duke_key_here

[Install]
WantedBy=multi-user.target
